# scripts/cron/com.neotrade.five_flags_pool.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neotrade.five_flags_pool</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/mac/NeoTrade2/scripts/cron/five_flags_pool_task.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mac/NeoTrade2</string>

    <key>StandardOutPath</key>
    <string>/Users/mac/NeoTrade2/logs/five_flags_pool.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mac/NeoTrade2/logs/five_flags_pool_error.log</string>

    <key>StartCalendarInterval</key>
    <integer>86400</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/mac/NeoTrade2</string>
    </dict>

    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>
